# deploy/pbs-nightly.timer
[Unit]
Description=Run PBS nightly at 03:00 (VM local time = JST)

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
